# files/common/usr/lib/systemd/user/notmuch-fetch.service
[Unit]
Description=Notmuch email fetcher and indexer

[Service]
Environment=NOTMUCH_CONFIG=%E/notmuch/notmuchrc
ExecStart=/usr/bin/notmuch new
